# config/crontab/production
50 23 * * * sh -c "cd /home/hlx/www/jobware/current && /usr/local/rvm/gems/ree-1.8.7-2010.02/bin/rake backup:database RAILS_ENV=production"
10 5 * * * sh -c "cd /home/hlx/www/jobware/current && /usr/local/rvm/gems/ree-1.8.7-2010.02/bin/rake -s sitemap:refresh RAILS_ENV=production"